# etc/supervisor.conf
[program:kandu]
command=/home/marlin/.virtualenvs/kandu/bin/python manage.py runserver 0.0.0.0:9090
directory=/home/marlin/apps/kandu
stopasgroup=true
redirect_stderr=true
